# etc/systemd/system/alsa-audio-server.service
[Unit]
Description=Stream audio from ALSA to network

[Service]
ExecStart=/usr/local/bin/alsa-audio-server

[Install]
WantedBy=multi-user.target
